# deploy/clientst0r-auto-update.service
[Unit]
Description=Client St0r Auto-Update Service
After=network.target

[Service]
Type=oneshot
User=clientst0r
Group=clientst0r
WorkingDirectory=/opt/clientst0r
ExecStart=/opt/clientst0r/scripts/auto_update.sh
StandardOutput=append:/var/log/clientst0r/auto-update.log
StandardError=append:/var/log/clientst0r/auto-update.log

# Load deployment settings (AUTO_UPDATE_ENABLED, etc.). Leading `-` = optional file.
EnvironmentFile=-/etc/clientst0r/.env
Environment="PATH=/opt/clientst0r/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# Security
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
