# deploy/clientst0r-psa-sync.service
[Unit]
Description=Client St0r PSA Sync Service
After=network.target

[Service]
Type=oneshot
User=clientst0r
Group=clientst0r
WorkingDirectory=/opt/clientst0r
Environment="PATH=/opt/clientst0r/venv/bin"
EnvironmentFile=-/etc/clientst0r/.env
ExecStart=/opt/clientst0r/venv/bin/python manage.py sync_psa
StandardOutput=journal
StandardError=journal
