# deploy/clientst0r-breach-scan.service
[Unit]
Description=Client St0r Password Breach Scan
After=network.target

[Service]
Type=oneshot
User=clientst0r
Group=clientst0r
WorkingDirectory=/opt/clientst0r
Environment="PATH=/opt/clientst0r/venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=-/etc/clientst0r/.env
ExecStart=/opt/clientst0r/venv/bin/python /opt/clientst0r/manage.py check_password_breaches
StandardOutput=journal
StandardError=journal
