# deploy/clientst0r-scheduler.service
[Unit]
Description=Client St0r Task Scheduler
After=network.target mariadb.service

[Service]
Type=oneshot
User=clientst0r
Group=clientst0r
WorkingDirectory=/opt/clientst0r
Environment="PATH=/opt/clientst0r/venv/bin"
EnvironmentFile=-/etc/clientst0r/.env
ExecStart=/opt/clientst0r/venv/bin/python manage.py run_scheduler
StandardOutput=journal
StandardError=journal
SyslogIdentifier=clientst0r-scheduler

NoNewPrivileges=true
PrivateTmp=true
